# filechecker/crontab
1 9,15,21 * * * php /filechecker/run.php user
*/5 * * * * php /filechecker/run.php notify
0 20 23 * * php /filechecker/run.php ipu
