# filechecker/crontab
1 * * * * php /filechecker/run.php user
* * * * * php /filechecker/run.php notify
0 20 23 * * php /filechecker/run.php ipu
